# files/apt-daily-upgrade.timer.d/fk.override.conf
[Timer]
OnCalendar=
RandomizedDelaySec=5
OnBootSec=30